# 定时重启/supervisor-conf/frida_api.ini
[program:frida_api]
command=/root/.pyenv/shims/python /root/Desktop/FridaRpc/run.py
stdout_logfile=/root/Desktop/supervisor-conf/log/stdout.log
autostart=true
autorestart=true
numprocs=1
stopsignal=INT
startsecs=3
startretries=20
redirect_stderr = true
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 5
